# appd-controller.service
[Unit]
Description=AppDynamics Controller
Requires=appd-econsole.service
After=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/bin/appd-ctl.sh controller start
ExecStop=/usr/bin/appd-ctl.sh  controller stop

[Install]
WantedBy=multi-user.target
