# appd-events-service.service
[Unit]
Description=AppDynamics Events Service
After=appd-econsole.service
Before=appd-econsole.service

[Service]
Type=simple
ExecStart=/usr/bin/appd-ctl.sh events_service start
ExecStop=/usr/bin/appd-ctl.sh  events_service stop
RemainAfterExit=yes


[Install]
WantedBy=appd-events-service.service
