# appd-events-service.service
[Unit]
Description=AppDynamics Events Service
Requires=appd-econsole.service
After=appd-controller.service
Before=appd-controller.service

[Service]
Type=simple
ExecStart=/usr/bin/appd-ctl.sh events_service start
ExecStop=/usr/bin/appd-ctl.sh  events_service stop
RemainAfterExit=yes


[Install]
WantedBy=multi-user.target
